# ansible/roles/influga/files/influga.conf
[program:influga]
command=/usr/local/bin/influga start -c /etc/influga/config.json
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
